# provider-dashboard-rs.service
[Unit]
Description=Provider Dashboard Rust Backend
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/provider_tracking
ExecStart=/home/user/provider_tracking/provider_tracker
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="RUST_LOG=info"

[Install]
WantedBy=multi-user.target
